# firmware/systemd/algoro-admin.service
[Unit]
Description=algoro Admin Web Interface
After=network.target unbound.service

[Service]
ExecStart=/usr/local/bin/uvicorn algoro.admin.app:app --host 0.0.0.0 --port 80
WorkingDirectory=/opt/algoro/firmware
Environment="ALGORO_DB_PATH=/var/lib/algoro/algoro.db"
Environment="ALGORO_BLOCKLIST_PATH=/var/lib/algoro/active_blocklist.txt"
Environment="ALGORO_UNBOUND_CONF=/etc/unbound/unbound.conf.d/algoro.conf"
Environment="ALGORO_DNS_TEMPLATE_DIR=/opt/algoro/firmware/dns"
Environment="ALGORO_PRESET_DIR=/opt/algoro/firmware/blocklists"
Restart=always
User=root

[Install]
WantedBy=multi-user.target
