# firmware/systemd/slopstop-admin.service
[Unit]
Description=slopstop Admin Web Interface
After=network.target unbound.service

[Service]
ExecStart=/usr/local/bin/uvicorn slopstop.admin.app:app --host 0.0.0.0 --port 80
WorkingDirectory=/opt/slopstop/firmware
Environment="SLOPSTOP_DB_PATH=/var/lib/slopstop/slopstop.db"
Environment="SLOPSTOP_BLOCKLIST_PATH=/var/lib/slopstop/active_blocklist.txt"
Environment="SLOPSTOP_UNBOUND_CONF=/etc/unbound/unbound.conf.d/slopstop.conf"
Environment="SLOPSTOP_DNS_TEMPLATE_DIR=/opt/slopstop/firmware/dns"
Environment="SLOPSTOP_PRESET_DIR=/opt/slopstop/firmware/blocklists"
Restart=always
User=root

[Install]
WantedBy=multi-user.target
